# scripts/piroutergui.service
[Unit]
Description=PiRouterGUI (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
ExecStart=__WORKDIR__/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

# Full-mode managed config targets
Environment=PRG_DNSMASQ_MANAGED_PATH=/etc/dnsmasq.d/piroutergui-static.conf
Environment=PRG_NFT_MANAGED_PATH=/etc/nftables.d/piroutergui-blocklist.nft
Environment="PRG_DNSMASQ_RELOAD_CMD=systemctl reload dnsmasq"
Environment="PRG_NFT_APPLY_CMD=nft -f /etc/nftables.d/piroutergui-blocklist.nft"
Environment="PRG_AUTH_ENABLED=true"
Environment="PRG_AUTH_USERNAME=admin"
Environment="PRG_AUTH_PASSWORD=change-me"
Environment="PRG_AUTH_SECRET=change-this-secret"

[Install]
WantedBy=multi-user.target
